# defaultscr/crontab
* * * * * root /usr/local/sbin/usystem_adduser
